# ah_background.service
[Unit]
Description=Authors haven service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/bin/bash /home/ubuntu/start_ah.sh

[Install]
WantedBy=multi-user.target
